# luci-theme-netgearrainbow/root/etc/uci-defaults/30_luci-theme-netgearrainbow
#!/bin/sh
uci batch <<-EOF
	set luci.themes.NetgearRainbow=/luci-static/netgearrainbow
	set luci.main.mediaurlbase=/luci-static/netgearrainbow
	commit luci
EOF
exit 0
